# deploy/telegram-student-bot.service
[Unit]
Description=Telegram Student Verification Bot
After=network.target

[Service]
Type=simple
User=student-bot
Group=student-bot
WorkingDirectory=/opt/apps/student-verification-bot
EnvironmentFile=/opt/apps/student-verification-bot/.env
ExecStart=/opt/apps/student-verification-bot/.venv/bin/python /opt/apps/student-verification-bot/app.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
